# assets/system/smartpixeldashboard.service.template
[Unit]
Description=Smart Pixel Dashboard
After=bluetooth.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__WORKDIR__
ExecStart=__WORKDIR__/.venv/bin/python __WORKDIR__/startup.py
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
