# assets/system/smartpixelpanel.service.template
[Unit]
Description=Smart Pixel Panel Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__WORKDIR__
ExecStart=/usr/bin/python3 __WORKDIR__/startup.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
